# docker/supervisor.ini
[program:caddyserver]
process_name=Caddyserver
command=caddy -conf /Caddyfile -log stdout
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout


[program:php-fpm]
process_name=php-fpm7
command=php-fpm7
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout

[program:cron]
process_name=cron
command=crond -f
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
